# config/includes.chroot/etc/xdg/autostart/cyrethonion-mate.desktop
[Desktop Entry]
Type=Application
Name=Cyrethonion
Exec=cyrethonion-mate
OnlyShowIn=MATE;
StartupNotify=false
Terminal=false
Categories=System;Security;Network;
